# etc/supervisor/conf.d/courier-imap-ssl.conf
[program:courier-imap-ssl]
command=/usr/sbin/couriertcpd -address=0 -maxprocs=40 -maxperip=20 -nodnslookup -noidentlookup 993 /usr/lib/courier/courier/imaplogin /usr/bin/imapd Maildir
priority=4
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0